# files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
# disable login prompt on serial port
#::askconsole:/usr/libexec/login.sh

